# configs/crontab
@daily /usr/bin/python /var/www/mailinator/manage.py cleaner
